# install/provision/gunicorn.conf
[program:echo-gunicorn]
environment=DJANGO_SETTINGS_MODULE='echo.settings.development'
command = /home/vagrant/venv/bin/gunicorn --log-level debug --pythonpath /home/vagrant/echo echo.wsgi:application --user="vagrant" --bind 0.0.0.0:8000 --workers=3 --reload
user = vagrant
directory = /home/vagrant
autostart = true
autorestart = false
stdout_logfile=/var/log/supervisor/echo-gunicorn.log
redirect_stderr=true